# backend/backend.service
# =============================================================
#  systemd Unit File — Backend API (gunicorn)
# =============================================================
# Copy to /etc/systemd/system/backend.service
# Then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable backend
#   sudo systemctl start backend
#   sudo systemctl status backend
#   journalctl -u backend -f        ← live logs
# =============================================================

[Unit]
Description=3-Tier App — Python Flask Backend API
After=network.target

[Service]
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/app/backend

# gunicorn: 4 worker processes, binding on all interfaces port 5000
ExecStart=/home/ec2-user/app/backend/venv/bin/gunicorn \
          --workers 4 \
          --bind 0.0.0.0:5000 \
          --timeout 120 \
          --access-logfile - \
          --error-logfile - \
          app:app

Restart=on-failure
RestartSec=5s

EnvironmentFile=/home/ec2-user/app/backend/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
