# worker/worker.service
# =============================================================
#  systemd Unit File — Background Worker
# =============================================================
# Copy this file to /etc/systemd/system/worker.service
# Then run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable worker
#   sudo systemctl start worker
#   sudo systemctl status worker
# =============================================================

[Unit]
Description=3-Tier App — Background Worker
After=network.target

[Service]
# --- User to run the process as (never run as root in production) ------------
User=ec2-user
Group=ec2-user

# --- Working directory (where your worker code lives) -----------------------
WorkingDirectory=/home/ec2-user/app/worker

# --- The command to run ------------------------------------------------------
ExecStart=/home/ec2-user/app/worker/venv/bin/python worker.py

# --- Restart on failure (with a 5-second delay) -----------------------------
Restart=on-failure
RestartSec=5s

# --- Load environment variables from .env file ------------------------------
EnvironmentFile=/home/ec2-user/app/worker/.env

# --- stdout/stderr go to the systemd journal (view with: journalctl -u worker)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
